# rake_test/rake-sensor-fixed.service
[Unit]
Description=Soil Monitor - Continuous Environmental Sensor Logger
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=maggi
Group=maggi
WorkingDirectory=/home/maggi/tilden_test/rake_test
Environment=HOME=/home/maggi
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=SDL_VIDEODRIVER=kmsdrm
Environment=SDL_FBDEV=/dev/fb0
Environment=SDL_NOMOUSE=1

# The start script will wait for graphics subsystem to be ready
ExecStart=/home/maggi/tilden_test/rake_test/start_forest_rings.sh

# Only restart on failure, not on clean exit
Restart=on-failure
RestartSec=10
StartLimitBurst=3
StartLimitIntervalSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=soil-monitor

[Install]
WantedBy=multi-user.target
