# rake_test/rake-sensor.service
[Unit]
Description=Rake Env + GPS Logger GUI
After=network-online.target graphical.target
Wants=network-online.target graphical.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/tilden_test/rake_test
# If using virtualenv uncomment next line and comment the system python line
# ExecStart=/home/pi/tilden_test/rake_test/.venv/bin/python /home/pi/tilden_test/rake_test/main.py
ExecStart=/usr/bin/python3 /home/pi/tilden_test/rake_test/main.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
# --- Environment for GUI display ---
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=SDL_AUDIODRIVER=dummy

[Install]
WantedBy=graphical.target
